# files/etc/systemd/system/zabbix-server.service.d/override.conf
# Managed by ansible (mit.zabbix-server.package)
# DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
PartOf=mariadb.service
[Install]
WantedBy=mariadb.service

